# udev/rules.d/99-bluetooth-dongle.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="0bda", ATTRS{idProduct}=="a729", SYMLINK+="bluetooth-dongle", RUN+="/usr/bin/systemctl restart bluetooth"
