# systemd/system/ir-listener.service
[Unit]
Description=IR Listener Service

[Service]
ExecStart=/path/to/your/compiled/program

[Install]
WantedBy=multi-user.target
